# docker/mastodon/periodic/hourly/backup
rsync -avz /mastodon/public/assets /backup/assets/
rsync -avz /mastodon/public/system /backup/system/