# docker/postgres/periodic/daily/postgres-dump
#!/bin/bash
ts=$(date +%d-%m-%Y"_"%H_%M_%S)
pg_dumpall -c -U postgres | gzip > /tmp/$ts.gz
cp /tmp/$ts.gz /backup/